# vendor/scripts/shinken-systemd.service
[Unit]
Description=Shinken daemon

[Service]
Type=oneshot
ExecStart=/etc/init.d/shinken start
ExecStop=/etc/init.d/shinken stop

[Install]
WantedBy=multi-user.target