# hokuto/etc/systemd/system/hokuto.service
[Unit]
Description=Hokuto daemon

[Service]
Type=simple
ExecStart=/usr/bin/hokuto.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target