# nanto/etc/systemd/system/nanto.service
[Unit]
Description=Nanto daemon

[Service]
Type=simple
ExecStart=/usr/bin/nanto.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target